# config/supervisor.conf
[unix_http_server]
file=/tmp/supervisor.sock ; (the path to the socket file)


[inet_http_server] ; inet (TCP) server disabled by default
port=*:9010 ; (ip_address:port specifier, *:port for all iface)


[supervisord]
logfile=/var/praekelt/vumi/supervisor.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10 ; (num of main logfile rotation backups;default 10)
loglevel=info ; (log level;default info; others: debug,warn,trace)
pidfile=/var/praekelt/vumi/supervisor.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false ; (start in foreground if true;default false)
minfds=1024 ; (min. avail startup file descriptors;default 1024)
minprocs=200 ; (min. avail process descriptors;default 200)
environment=TZ="Africa/Johannesburg"


[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:vumi]
priority=999
directory=/var/praekelt/vumi
command=twistd -n 
  --pidfile=/var/praekelt/vumi/tmp/pids/%(program_name)s_%(process_num)s.pid
  start_worker 
  --config=/var/praekelt/vumi/config.yaml 
  --worker-class=vumi.multiworker.MultiWorker
  --vhost=/develop      
process_name=%(program_name)s_%(process_num)s
numprocs=1
numprocs_start=1
autostart=true
autorestart=true

stdout_logfile=/var/praekelt/vumi/logs/vumi.log
redirect_stderr = true
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10


[program:webpy]
priority=998
directory=/var/thoughtworks/web_sms_app
command=python shoprite_webapp.py
process_name=%(program_name)s_%(process_num)s
numprocs=1
numprocs_start=1
autostart=true
autorestart=true

stdout_logfile=/var/thoughtworks/web_sms_app/logs/webpy.log
redirect_stderr = true
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10


[program:mongodb]
priority=997
directory=/var/thoughtworks/mongodb
command=service mongodb restart
process_name=%(program_name)s_%(process_num)s
numprocs=1
numprocs_start=1
autostart=true
autorestart=true

stdout_logfile=/var/thoughtworks/mongodb/logs/mongodb.log
redirect_stderr = true
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10